# audio-daemon/service-templates/teamsigma-audiodaemon.service
# REFERENCE template (human-readable). The authoritative unit is generated at
# runtime by audio_daemon/platforms/linux.py with the correct ExecStart.
# Installed to: ~/.config/systemd/user/teamsigma-audiodaemon.service
# Token: __EXEC_START__  (e.g. "/path/to/python -m audio_daemon --mode=daemon")
[Unit]
Description=TeamSigma Audio Daemon
After=default.target sound.target
Wants=sound.target

[Service]
Type=simple
ExecStart=__EXEC_START__
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3s
Environment=AUDIO_DAEMON_PORT=8770
# PipeWire/PulseAudio reachability from a --user service:
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus
Environment=PULSE_RUNTIME_PATH=/run/user/%U/pulse
# Lightweight hardening (loose enough to write logs/state under $HOME):
MemoryMax=150M
Nice=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=true
SyslogIdentifier=teamsigma-audiodaemon
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
